# hhvm.sv.conf
[program:hhvm]
command=/bin/bash -c "hhvm -m server -p 8080 -d hhvm.server.type=proxygen -d hhvm.server.port=8080 -d hhvm.server.source_root=/srv/www/"
redirect_stderr=true
